# lib/systemd/system/virtualbox-guest-utils.service.d/40_vbox-disable-timesync.conf
## Copyright (C) 2016 - 2018 ENCRYPTED SUPPORT LP <adrelanos@riseup.net>
## See the file COPYING for copying conditions.

[Service]
ExecStart=
ExecStart=/etc/init.d/virtualbox-guest-utils.anondist start

ExecStop=
ExecStop=/etc/init.d/virtualbox-guest-utils.anondist stop
